# data/debian/sonic-host-services-data.gnoi-shutdown.service
[Unit]
Description=gNOI based DPU Graceful Shutdown Daemon
Requires=database.service
Wants=network-online.target
After=network-online.target database.service

[Service]
Type=simple
ExecCondition=/usr/bin/python3 /usr/local/bin/check_platform.py
ExecStartPre=/bin/bash /usr/local/bin/wait-for-sonic-core.sh
ExecStart=/usr/bin/python3 /usr/local/bin/gnoi_shutdown_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
